# files_for_server/start_bots.conf
[program:bot_alphabtc]
command=/home/.env/bin/list-bot/ALPHA_BTC.sh
autostart=true
autorestart=true
stderr_logfile=/home/log/bot_alphabtc.err.log
stdout_logfile=/home/log/bot_alphabtc.out.log

